# linuxRMS/debian/rmsgw.cron.d
#
# Regular cron jobs for the rmsgw package
#
21,51 * * * *   root   /usr/bin/rmsgw_aci >/dev/null 2>&1
